# cron.d/checkup
*/10 * * * * root /home/pi/pi-distance/checkup
